# docker/development/etc/s6-overlay/s6-rc.d/init-setup/up
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv
cd /var/www/html
foreground {
    mkdir
    -p
    /var/www/html/storage/app/ssh/keys
}
foreground {
    mkdir
    -p
    /var/www/html/storage/app/ssh/mux
}
foreground {
    chmod
    -R
    700
    /var/www/html/storage/app/ssh
}
foreground {
    chown
    -R
    www-data:www-data
    /var/www/html/storage/app/ssh
}
foreground {
    mkdir
    -p
    /var/www/html/storage/logs
}
foreground {
    touch
    /var/www/html/storage/logs/laravel.log
}
foreground {
    chown
    -R
    www-data:www-data
    /var/www/html/storage/logs
}
foreground {
    chmod
    -R
    755
    /var/www/html/storage/logs
}
foreground {
    chmod
    664
    /var/www/html/storage/logs/laravel.log
}
foreground {
    composer
    install
}
foreground {
    php
    artisan
    migrate
    --step
}
foreground {
    php
    artisan
    dev
    --init
}

