# docker/production/etc/s6-overlay/s6-rc.d/init-script/up
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv
cd /var/www/html
foreground {
    mkdir
    -p
    /var/www/html/storage/logs
}
foreground {
    touch
    /var/www/html/storage/logs/laravel.log
}
foreground {
    chown
    -R
    www-data:www-data
    /var/www/html/storage/logs
}
foreground {
    chmod
    -R
    755
    /var/www/html/storage/logs
}
foreground {
    chmod
    664
    /var/www/html/storage/logs/laravel.log
}
foreground {
    php
    artisan
    app:init
}


